# deploy/caddy/systemd/caddy.service.d/override.conf
[Service]
Restart=always
RestartSec=5s
LimitNOFILE=1048576
